# init/gpu-tune
#!/sbin/openrc-run
description="AMD GPU tuning"
start() {
    bash /usr/local/bin/gpu-tune.sh
}
